# src/etc/services.d/nmbd/finish
#! /bin/sh

echo "nmbd is stopped"
s6-svscanctl -t /var/run/s6/services
