# src/etc/services.d/nmbd/run
#!/bin/sh

echo "nmbd is starting"
exec nmbd -i
