# src/etc/services.d/smbd/finish
#! /bin/sh

echo "smbd is stopped"
s6-svscanctl -t /var/run/s6/services
