# src/etc/services.d/smbd/run
#!/bin/sh

echo "smbd is starting"
exec smbd -i
